# image/src/modules/waypoint/filesystem/root/etc/systemd/system/waypoint-m17gateway.service
[Unit]
Description=Waypoint M17Gateway
After=waypoint-mmdvm.service

[Service]
# Runs the waypoint-stack package binary (/usr/bin) against the INI waypointd
# renders on Apply. Managed by waypointd: enabled/started when its mode is
# configured, health-gated on stack updates (RFC-0014). Not enabled at image
# build — a fresh, unclaimed node runs no gateways until it is configured.
ExecStart=/usr/bin/M17Gateway /home/pi-star/waypoint/etc/M17Gateway.ini
Restart=on-failure
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
